# pi_hub_/syswatch-hub.service
# Servicio systemd para la Raspberry Pi Zero 2W
# Copiar a: /etc/systemd/system/syswatch-hub.service
# Activar:
#   sudo systemctl daemon-reload
#   sudo systemctl enable syswatch-hub
#   sudo systemctl start syswatch-hub

[Unit]
Description=SysWatch Hub — Monitor de equipos
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/syswatch/hub
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=10
StandardOutput=append:/home/pi/syswatch/hub.log
StandardError=append:/home/pi/syswatch/hub.log

# Limitar recursos en la Pi Zero
MemoryMax=180M
CPUQuota=80%

# Variables de entorno — editar aquí o en archivo .env
Environment=SYSWATCH_API_KEY=clave-secreta-cambiar
Environment=SYSWATCH_PG_DSN=postgresql://user:pass@host/syswatch
Environment=SYSWATCH_PG_ON=true
Environment=SYSWATCH_LOG=INFO

[Install]
WantedBy=multi-user.target
